# flask_bbs.conf
[program:flask_bbs]
command=/usr/local/bin/gunicorn wsgi --bind localhost:3000 --workers 3 --worker-class gevent
directory=/var/www/flask_bbs
autostart=true
autorestart=true
